# files/system/usr/lib/systemd/user/obex.service.d/10-symlinks.conf
# Enable --symlinks flag for obexd to allow Bluetooth file receiving
# on atomic/immutable systems where /home is a symlink to /var/home.
# Without this flag, obexd refuses to follow symlinks for security reasons,
# causing "Operation not permitted" errors when receiving files.

[Service]
ExecStart=
ExecStart=/usr/libexec/bluetooth/obexd --symlinks
